# scripts/recent_files_monitor/cleanup_recent_activity.service
[Unit]
Description=Cleanup StoragePod activity log (entries older than 7 days)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/storagepod_cleanup.sh